# roles/common/files/84-thinkpad-dock.rules
SUBSYSTEM=="usb", ACTION=="add|remove", ENV{ID_VENDOR_ID}=="17ef", ENV{ID_MODEL_ID}=="1010", RUN+="/usr/local/bin/monitors.sh"
